# deploy/systemd/rosadmin@staging.service.d/override.conf
# Staging instance of rosadmin@.service.
[Service]
User=rosadmin-staging
Group=rosadmin-staging
# The origin the API expects on state-changing browser requests.
Environment=ROSADMIN_ORIGIN=https://rosadmin.staging.portlanddsa.org
# botonio's SSO endpoint: the socket rosadmin relays to, plus the public values it
# verifies an assertion against. The verifying key and guild id are filled from the
# encrypted secrets file when this unit is rendered; the socket path is fixed per stage.
Environment=BOTONIO_SSO_SOCKET_PATH=/run/botonio-staging/sso.sock
Environment=BOTONIO_SSO_PUBLIC_KEY=${botonio_sso_pubkey}
Environment=BOTONIO_SSO_GUILD_ID=${sso_guild_id}
# The Google Workspace mirror: the DWD impersonation subject is filled from the
# encrypted secrets file; the key itself arrives as a credential below, and
# CREDENTIALS_FILE points the app at it (%d is the unit's credentials directory).
# Staging expects example-domain member addresses (test personas), so their
# sync skips log quietly here.
Environment=ROSADMIN_GOOGLE_SUBJECT=${google_dwd_subject}
Environment=CREDENTIALS_FILE=%d/google-dwd-key
Environment=ROSADMIN_EXPECT_EXAMPLE_EMAILS=1
# The operator admin surface: a second uvicorn listener on a private socket in
# this unit's runtime directory. The DIRECTORY is the whole access control:
# uvicorn chmods its unix socket world-writable on bind and systemd's default
# directory mode is world-traversable, so without the 0750 here every local
# user on the co-tenant box could connect. With it, only this unit's user and
# root can reach the socket at all - possession is the authorization.
RuntimeDirectory=rosadmin-staging
RuntimeDirectoryMode=0750
Environment=ROSADMIN_ADMIN_SOCKET=/run/rosadmin-staging/admin.sock
# Staging's membership source is the in-process mock Solidarity Tech server:
# the service binds it to the base URL below at boot (loopback only), the
# roster pull reads it, and the admin socket's persona relay drives it. The
# persona map is only the boot roster - the relay mutates standing and
# leadership at runtime - and it is filled from the encrypted secrets file.
Environment=SOLIDARITY_TECH_MOCK=1
Environment=SOLIDARITY_TECH_BASE_URL=http://127.0.0.1:8001
Environment=SOLIDARITY_TECH_MOCK_PERSONAS=${st_mock_personas}
# Secrets delivered as host-key-encrypted systemd credentials - readable only by
# this unit's user through the credentials directory. 
# Provisioned from the encrypted secrets by the deploy tool; see the
# activation runbook. botonio_sso_bearer is the shared secret rosadmin presents to
# the SSO socket.
LoadCredentialEncrypted=audit-hmac-key:/etc/rosadmin/service/staging/audit-hmac-key.cred
LoadCredentialEncrypted=db_migration_password:/etc/rosadmin/service/staging/db_migration_password.cred
LoadCredentialEncrypted=botonio_sso_bearer:/etc/rosadmin/service/staging/botonio_sso_bearer.cred
LoadCredentialEncrypted=google-dwd-key:/etc/rosadmin/service/staging/google-dwd-key.cred
