# deploy/systemd/rosadmin@staging.service.d/override.conf
# Staging instance of rosadmin@.service.
[Service]
User=rosadmin-staging
Group=rosadmin-staging
# Staging serves the development surface: fake-login and the persona-backed
# stubs. The flag is half the gate - the other half is the rosadmin-devtools
# package, which only staging artifacts contain.
Environment=ROSADMIN_FAKE_LOGIN=1
# The origin the API expects on state-changing browser requests.
Environment=ROSADMIN_ORIGIN=https://rosadmin.staging.portlanddsa.org
