# deploy/systemd/rosadmin@staging.service.d/override.conf
# Staging instance of rosadmin@.service - serves only /api/healthz.
[Service]
User=rosadmin-staging
Group=rosadmin-staging
Environment=ROSADMIN_PORT=40643
