# deploy/systemd/rosadmin@staging.service.d/override.conf
# Staging instance of rosadmin@.service.
[Service]
User=rosadmin-staging
Group=rosadmin-staging
# The origin the API expects on state-changing browser requests.
Environment=ROSADMIN_ORIGIN=https://rosadmin.staging.portlanddsa.org
# botonio's SSO endpoint: the socket rosadmin relays to, plus the public values it
# verifies an assertion against. The verifying key and guild id are filled from the
# encrypted secrets file when this unit is rendered; the socket path is fixed per stage.
Environment=BOTONIO_SSO_SOCKET_PATH=/run/botonio-staging/sso.sock
Environment=BOTONIO_SSO_PUBLIC_KEY=${botonio_sso_pubkey}
Environment=BOTONIO_SSO_GUILD_ID=${sso_guild_id}
# Secrets delivered as host-key-encrypted systemd credentials - readable only by
# this unit's user through the credentials directory, never the process
# environment. Provisioned from the encrypted secrets by the deploy tool; see the
# activation runbook. botonio_sso_bearer is the shared secret rosadmin presents to
# the SSO socket.
LoadCredentialEncrypted=audit-hmac-key:/etc/rosadmin/service/staging/audit-hmac-key.cred
LoadCredentialEncrypted=db_migration_password:/etc/rosadmin/service/staging/db_migration_password.cred
LoadCredentialEncrypted=botonio_sso_bearer:/etc/rosadmin/service/staging/botonio_sso_bearer.cred
