# deploy/systemd/rosadmin@staging.service.d/override.conf
# Staging instance of rosadmin@.service.
[Service]
User=rosadmin-staging
Group=rosadmin-staging
# Staging serves the development surface: fake-login and the persona-backed
# stubs. The flag is half the gate - the other half is the rosadmin-devtools
# package, which only staging artifacts contain.
Environment=ROSADMIN_FAKE_LOGIN=1
# The origin the API expects on state-changing browser requests.
Environment=ROSADMIN_ORIGIN=https://rosadmin.staging.portlanddsa.org
# The audit HMAC key and the migration role's scram password, delivered as
# host-key-encrypted systemd credentials - readable only by this unit's user
# through $CREDENTIALS_DIRECTORY, never the process environment. Provisioned from
# the encrypted secrets by the deploy tool; see the activation runbook.
LoadCredentialEncrypted=audit-hmac-key:/etc/rosadmin/service/staging/audit-hmac-key.cred
LoadCredentialEncrypted=db_migration_password:/etc/rosadmin/service/staging/db_migration_password.cred
