# deploy/systemd/rosadmin@staging.service.d/override.conf
# Staging instance of rosadmin@.service - serves only /healthz.
[Service]
User=rosadmin-staging
Group=rosadmin-staging
# No LoadCredentialEncrypted= lines, no migrate ExecStartPre, and no sync timer here:
# the Google key, the Solidarity Tech token, the database, and the periodic reconcile
# all land when a database is activated. Staging needs none of them to answer /healthz.
