# deploy/systemd/rosadmin@production.service.d/override.conf
# Production instance of rosadmin@.service - the leader panel.
[Service]
User=rosadmin-production
Group=rosadmin-production
# Tenant and roster shared with the sweep through one rendered file, so the two
# cannot drift on which Workspace they impersonate. The Google subject lives there.
EnvironmentFile=/etc/rosadmin/service/%i/shared.env
# The origin the API expects on state-changing browser requests.
Environment=ROSADMIN_ORIGIN=https://rosadmin.portlanddsa.org
# botonio's production SSO endpoint: the socket rosadmin relays to, plus the public
# values it verifies an assertion against. The verifying key and guild id are filled
# from the encrypted secrets when this unit is rendered; the socket path is fixed.
Environment=BOTONIO_SSO_SOCKET_PATH=/run/botonio-production/sso.sock
Environment=BOTONIO_SSO_PUBLIC_KEY=${botonio_sso_pubkey}
Environment=BOTONIO_SSO_GUILD_ID=${sso_guild_id}
# The Google Workspace mirror key arrives as a credential below; CREDENTIALS_FILE
# points the app at its per-unit delivered path (%d is the unit's credentials
# directory). The impersonation subject lives in the shared env file above.
Environment=CREDENTIALS_FILE=%d/google-dwd-key
# The operator admin surface: a second uvicorn listener on a private socket in this
# unit's runtime directory. The DIRECTORY is the whole access control: only this
# unit's user and root can reach the 0750 directory, so possession of the socket is
# the authorization.
RuntimeDirectory=rosadmin-production
RuntimeDirectoryMode=0750
Environment=ROSADMIN_ADMIN_SOCKET=/run/rosadmin-production/admin.sock
# Production database names, NOT the staging defaults - the panel and its migrate
# ExecStartPre (on the base unit) would otherwise target staging. The runtime role
# reaches the database by peer over the socket; the migration role connects over TCP
# loopback with scram using the db_migration_password credential below.
Environment=ROSADMIN_DB_NAME=rosadmin_production
Environment=ROSADMIN_DB_USER=rosadmin_production_app
Environment=ROSADMIN_DB_MIGRATE_USER=rosadmin_production_migrate
# Secrets delivered host-key-encrypted by the `service` deploy component, whose dest
# dir is /etc/rosadmin/service/production. Readable only by this unit's user through
# the credentials directory. No Solidarity Tech token: the panel reads the local
# mirror, not the API.
LoadCredentialEncrypted=audit-hmac-key:/etc/rosadmin/service/production/audit-hmac-key.cred
LoadCredentialEncrypted=db_migration_password:/etc/rosadmin/service/production/db_migration_password.cred
LoadCredentialEncrypted=botonio_sso_bearer:/etc/rosadmin/service/production/botonio_sso_bearer.cred
LoadCredentialEncrypted=google-dwd-key:/etc/rosadmin/service/production/google-dwd-key.cred
