# deploy/systemd/rosadmin-sync@.service
[Unit]
Description=rosadmin roster pull and reconcile sweep (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/rosadmin/%i/current
# Boot-time integrity: verify the extracted tree against the manifest the build
# emitted, from WorkingDirectory where MANIFEST.sha256 lives. A mismatch aborts
# the run before it starts. The tree is read-only and the service user cannot
# write it, so this is tamper-evidence and drift detection, not a barrier against root.
ExecStartPre=/usr/bin/sha256sum -c --quiet --strict MANIFEST.sha256
# One pull-then-reconcile pass. The web service owns schema currency (it migrates
# on its own start), so this unit does not migrate. Per-instance User=/Group=,
# credentials, the shared env file, and the staging dry-run flag live in the
# matching rosadmin-sync@<instance>.service.d/override.conf drop-in.
ExecStart=/opt/rosadmin/%i/current/python/bin/python3 -m rosadmin sync run
Environment=PYTHONDONTWRITEBYTECODE=1
# The first population of an empty group can add a few thousand members; the
# steady state is seconds. A generous bound still stops a wedged run from hanging
# the timer forever.
TimeoutStartSec=1800

# --- sandbox (validate on the box: systemd-analyze security rosadmin-sync@staging) ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateUsers=true
ProtectProc=invisible
ProcSubset=pid
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RemoveIPC=true
UMask=0077
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
# The database is reached by peer over an AF_UNIX socket; AF_INET/AF_INET6 remain
# for the outbound Solidarity Tech and Google API calls.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# MemoryDenyWriteExecute is deliberately NOT set: cryptography's cffi callbacks
# (reachable via google-auth RS256 when the Google credential is loaded) may need
# W^X relaxed. Validate on the box before enabling.
