# deploy/systemd/rosadmin-sync@staging.service.d/override.conf
# Staging instance of rosadmin-sync@.service.
[Unit]
# On staging the Solidarity Tech source is the mock server bound inside the web
# service's own process, so order the sweep after it to avoid a connection
# refused right after a redeploy. Soft ordering only - a down web service must
# not block the timer from trying, and a failed run self-heals on the next
# window. Production reads the real Solidarity Tech API and has no such coupling,
# which is why this lives in the staging drop-in, not the base unit.
After=rosadmin@%i.service

[Service]
User=rosadmin-staging
Group=rosadmin-staging
# The tenant and roster this unit sees come from the shared env file, rendered
# from the same secrets as the web service, so the two cannot drift on which
# Workspace they impersonate or which Solidarity Tech roster they read.
EnvironmentFile=/etc/rosadmin/service/%i/shared.env
# The org-wide group the whole good-standing roster syncs into, filled from the
# encrypted secrets file - an org address kept out of the committed template for
# the same privacy reason as the impersonation subject.
Environment=ROSADMIN_MAIN_GROUP_EMAIL=${main_group_email}
# The Google Workspace mirror key arrives as a credential below; CREDENTIALS_FILE
# points the app at its per-unit delivered path (%d is this unit's credentials dir).
Environment=CREDENTIALS_FILE=%d/google-dwd-key
# Staging rehearses in dry-run - real Google reads, the real diff, recorded no-op
# writes - until the soak clears. Removing the flag is a deploy, never a hand edit
# on the box. The empty ExecStart= resets the oneshot's single start line before
# the replacement.
ExecStart=
ExecStart=/opt/rosadmin/%i/current/python/bin/python3 -m rosadmin sync run
# The sweep writes audit rows and mirrors to Google, so it loads only those two
# credentials - not the migration password (it does not migrate) or the SSO
# bearer (it authenticates no users).
LoadCredentialEncrypted=audit-hmac-key:/etc/rosadmin/service/staging/audit-hmac-key.cred
LoadCredentialEncrypted=google-dwd-key:/etc/rosadmin/service/staging/google-dwd-key.cred
