# deploy/systemd/rosadmin-sync@production.service.d/override.conf
# Production instance of rosadmin-sync@.service.
#
# STOPGAP: production runs the reconcile sweep alone - no web service yet. That
# makes this override self-contained where staging splits config across a shared
# env file the web service and the sweep both read, and where staging leaves
# migration to the web unit. When the production web service lands, extract the
# shared keys (the Google subject, the Solidarity Tech source, the database name)
# into a rendered shared.env like staging, move the migrate step back to the web
# unit, and fold the `sync` deploy component back into `service`. See the
# production sweep runbook.
[Service]
User=rosadmin-production
Group=rosadmin-production

# --- Solidarity Tech: the real API, no mock ---
# SOLIDARITY_TECH_MOCK is intentionally unset - production reads the real roster;
# SOLIDARITY_TECH_BASE_URL is unset too, so the client falls back to the public
# API. The token arrives as a credential below. This share of the shared key's
# request budget sums with the Discord bot's under Solidarity Tech's 60-per-30s
# limit; the bot mirrors a smaller roster and rarely contends, so the sweep takes
# the larger share.
Environment=SOLIDARITY_TECH_REQUEST_BUDGET=50

# --- Google Workspace mirror ---
# The real portlanddsa.org tenant admin the DWD service account impersonates,
# filled from the encrypted secrets at render time. The DWD key arrives as a
# credential; CREDENTIALS_FILE points the app at its per-unit delivered path
# (%d is this unit's credentials directory).
Environment=ROSADMIN_GOOGLE_SUBJECT=${google_dwd_subject}
Environment=CREDENTIALS_FILE=%d/google-dwd-key

# --- Self-provisioning group structure ---
# The org-wide group the whole good-standing roster syncs into, and the domain
# every minted group address takes. The main group's display name is used only
# when the sweep must create it; an existing group is adopted unchanged.
# ROSADMIN_MASS_CREATION_TRIPWIRE is left unset: its code default (10) is a safe
# steady-state cap, and the first bootstrap run ignores it by design - which is
# exactly why the first production run rehearses in dry-run (below).
Environment=ROSADMIN_MAIN_GROUP_EMAIL=${main_group_email}
Environment=ROSADMIN_MAIN_GROUP_NAME=${main_group_name}
Environment=ROSADMIN_GROUP_EMAIL_DOMAIN=portlanddsa.org

# --- Database: the production names, NOT the staging defaults ---
# Both the runtime DSN and the migrate URI default to the staging database and
# roles when unset, so production MUST name its own or the sweep would silently
# reconcile against staging's database. The runtime role reaches the database by
# peer over the socket; the migration role connects over TCP loopback with scram.
Environment=ROSADMIN_DB_NAME=rosadmin_production
Environment=ROSADMIN_DB_USER=rosadmin_production_app
Environment=ROSADMIN_DB_MIGRATE_USER=rosadmin_production_migrate

# --- Schema currency: the sweep migrates, because nothing else does ---
# STOPGAP: on staging the web service owns migrations through its own
# ExecStartPre; with no production web service, the sweep applies pending
# migrations before it runs. Drop-in ExecStartPre lines accumulate, so this runs
# after the base unit's manifest check and before the sweep. Uses the migration
# role over TCP loopback with the db_migration_password credential loaded below.
ExecStartPre=/opt/rosadmin/%i/current/python/bin/python3 -m rosadmin migrate

# --- The sweep itself: dry-run until the production soak clears ---
# Production rehearses in dry-run first - real Solidarity Tech and Google reads,
# the real diff and provisioning plan, recorded no-op writes - so the bootstrap
# run's group creations and the whole-roster main-group population can be reviewed
# in the journal before anything is created for real. Arming is removing
# --dry-run and redeploying, never a hand edit on the box. The empty ExecStart=
# resets the oneshot's single start line before the replacement.
ExecStart=
ExecStart=/opt/rosadmin/%i/current/python/bin/python3 -m rosadmin sync run --dry-run

# --- Credentials ---
# Delivered host-key-encrypted by the `sync` deploy component, whose dest dir is
# /etc/rosadmin/sync/production. Readable only by this unit's user through the
# credentials directory: the audit key, the DWD mirror key, the migration
# password, and the Solidarity Tech token. No SSO bearer - the sweep logs no one
# in.
LoadCredentialEncrypted=audit-hmac-key:/etc/rosadmin/sync/production/audit-hmac-key.cred
LoadCredentialEncrypted=google-dwd-key:/etc/rosadmin/sync/production/google-dwd-key.cred
LoadCredentialEncrypted=db_migration_password:/etc/rosadmin/sync/production/db_migration_password.cred
LoadCredentialEncrypted=solidarity_tech_token:/etc/rosadmin/sync/production/solidarity_tech_token.cred
