# deploy/systemd/rosadmin-sync@production.service.d/override.conf
# Production instance of rosadmin-sync@.service - the reconcile sweep.
[Service]
User=rosadmin-production
Group=rosadmin-production
# Tenant and roster from the shared env file, rendered from the same secrets as the
# panel, so the two cannot drift on which Workspace they impersonate or which
# Solidarity Tech roster they read. The real API, the request budget, and the Google
# subject all live there.
EnvironmentFile=/etc/rosadmin/service/%i/shared.env
# The org-wide group the whole good-standing roster syncs into, filled from the
# encrypted secrets file - an org address kept out of the committed template.
Environment=ROSADMIN_MAIN_GROUP_EMAIL=${main_group_email}
# The domain every group address the sweep mints takes. Non-secret and stage-fixed.
Environment=ROSADMIN_GROUP_EMAIL_DOMAIN=portlanddsa.org
# The main group's display name, paired with its address in the secrets file. Used
# only if the sweep must create the main group; an existing one is adopted unchanged.
# ROSADMIN_MASS_CREATION_TRIPWIRE is intentionally unset - its code default (10) is a
# safe steady-state cap, and the first bootstrap run ignored it by design.
Environment=ROSADMIN_MAIN_GROUP_NAME=${main_group_name}
# The Google Workspace mirror key arrives as a credential below; CREDENTIALS_FILE
# points the app at its per-unit delivered path (%d is this unit's credentials dir).
Environment=CREDENTIALS_FILE=%d/google-dwd-key
# Production names its own database (staging is the default). The runtime role reaches
# the database by peer over the socket.
Environment=ROSADMIN_DB_NAME=rosadmin_production
Environment=ROSADMIN_DB_USER=rosadmin_production_app
Environment=ROSADMIN_DB_MIGRATE_USER=rosadmin_production_migrate
# Armed and live: the sweep reconciles real Google groups on its 4-hourly timer. The
# empty ExecStart= resets the oneshot's single start line before the replacement.
ExecStart=
ExecStart=/opt/rosadmin/%i/current/python/bin/python3 -m rosadmin sync run
# The sweep writes audit rows, mirrors to Google, and reads the real Solidarity Tech
# roster - no migration password (the panel owns migration through its own
# ExecStartPre) and no SSO bearer (it authenticates no users). Delivered by the
# `service` deploy component to /etc/rosadmin/service/production.
LoadCredentialEncrypted=audit-hmac-key:/etc/rosadmin/service/production/audit-hmac-key.cred
LoadCredentialEncrypted=google-dwd-key:/etc/rosadmin/service/production/google-dwd-key.cred
LoadCredentialEncrypted=solidarity_tech_token:/etc/rosadmin/service/production/solidarity_tech_token.cred
